# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

# Can't use setuid in CentOS with Upstart 0.6.5
#setuid urbushey
chdir /home/urbushey/sites/SITENAME/source

exec ../virtualenv/bin/gunicorn --bind unix:/tmp/SITENAME.socket bullets.wsgi:application
